# deploy/idilica-api.service
# systemd unit — supervisa el binario Go (corre como root: es el único
# usuario del droplet).
#
# Instalación (una vez):
#   scp deploy/idilica-api.service root@<host>:/etc/systemd/system/
#   ssh root@<host> 'systemctl daemon-reload && systemctl enable idilica-api'
#
# Operación diaria:
#   systemctl status idilica-api
#   systemctl restart idilica-api
#   journalctl -u idilica-api -f
[Unit]
Description=Idilica API (Go)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
# El binario lee el .env de su WorkingDirectory
WorkingDirectory=/home/ellebkey/apps/idilica/backend
ExecStart=/home/ellebkey/apps/idilica/backend/api
Restart=always
RestartSec=3
PrivateTmp=true

[Install]
WantedBy=multi-user.target
